# scripts/nexus-pillar.service
# systemd unit for the NEXUS Pillar Server.
#
# Install:
#   sudo cp scripts/nexus-pillar.service /etc/systemd/system/
#   sudo cp build/app/pillard /usr/local/bin/
#   sudo useradd --system --home /var/lib/nexus --shell /usr/sbin/nologin nexus
#   sudo install -d -o nexus -g nexus -m 0700 /var/lib/nexus
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nexus-pillar.service
#
# Inspect:
#   systemctl status nexus-pillar
#   journalctl -u nexus-pillar -f
#   sudo systemctl kill --signal=SIGUSR1 nexus-pillar   # dump stats
#   sudo systemctl reload nexus-pillar                   # re-announce (SIGHUP)
#
# Admin UDS (non-root if your user is in the `nexus` group):
#   echo reload   | socat - UNIX-CONNECT:/run/nexus/pillard.sock
#   echo stats    | socat - UNIX-CONNECT:/run/nexus/pillard.sock
#   echo shutdown | socat - UNIX-CONNECT:/run/nexus/pillard.sock

[Unit]
Description=NEXUS Pillar Server (public mesh relay)
Documentation=https://github.com/anthropics/nexus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/var/lib/nexus

# Foreground (-f) is required under Type=simple so systemd tracks the PID.
# Adjust port/peers/flags as needed.
ExecStart=/usr/local/bin/pillard -f -p 4242 -i /var/lib/nexus/pillar.identity -u /run/nexus/pillard.sock
ExecReload=/bin/kill -HUP $MAINPID

# systemd creates /run/nexus at start (mode 0755, owned by User+Group) and
# removes it on stop. The admin UDS lives there so operators in the `nexus`
# group can issue reload/stats without sudo.
RuntimeDirectory=nexus
RuntimeDirectoryMode=0750

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/nexus
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Listen on :4242 by default; no privileged port binding needed.
# Uncomment if you want to bind <1024:
#AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource caps
LimitNOFILE=65536
TasksMax=256

# Logs go to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-pillar

[Install]
WantedBy=multi-user.target
